# supervisor.conf
[program:sshd]
user=root
command=/usr/sbin/sshd
redirect_stderr=true
