# supervisor.conf
[program:sshd]
user=root
command=/start_sshd.sh
redirect_stderr=true
